# SHLiC/bin/shlic-controller.service
[Unit]
Description=gunicorn for shlic-controller
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/games:/usr/games"
ExecStart=/usr/bin/gunicorn --workers 3 --bind unix:shlic-controller.sock -m 000 runController:controller.engine.app

[Install]
WantedBy=multi-user.target